# crontab
# min hour day month weekday command
0     1    *   *     *       . venv/bin/activate && flask train